# shell/firewalld-ui.service
[Unit]
Description=Firewalld UI Application
After=network.target firewalld.service

[Service]
Type=forking
User=root
WorkingDirectory=/workspaces/Firewalld-UI
ExecStart=__NPM_EXEC_PATH__ start
PIDFile=%H/run/egg-server.pid
Restart=always
RestartSec=10
Environment=NODE_ENV=production

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target